# etc/init.d/nvidia-suspend
#!/sbin/openrc-run

description="NVIDIA system suspend actions"

start() {
    ebegin "Running NVIDIA system suspend actions"
    /usr/bin/logger -t suspend -s "nvidia-suspend.service"
    /usr/bin/nvidia-sleep.sh "suspend"
    eend $?
}
