# etc/init.d/nvidia-hibernate
#!/sbin/openrc-run

description="NVIDIA system hibernate actions"

start() {
    ebegin "Running NVIDIA system hibernate actions"
    /usr/bin/logger -t hibernate -s "nvidia-hibernate.service"
    /usr/bin/nvidia-sleep.sh "hibernate"
    eend $?
}
