# etc/init.d/nvidia-powerd
#!/sbin/openrc-run

description="nvidia-powerd service"
command="/usr/bin/nvidia-powerd"
command_background="yes"

depend() {
    need localmount
    after bootmisc
    use logger
}
