# etc/init.d/nvidia-powerd
#!/sbin/openrc-run

description="Manages dynamic power states for NVIDIA GPUs"
command="/usr/bin/nvidia-powerd"
supervisor=supervise-daemon

depend() {
	need modules
	after localmount
}
