# etc/init.d/nvidia-persistenced-openrc
#!/sbin/openrc-run

description="NVIDIA Persistence Daemon"
command="/usr/bin/nvidia-persistenced"
command_args="--user nvpd"
pidfile="/run/nvidia-persistenced/nvidia-persistenced.pid"

depend() {
	need modules
	after localmount
	use logger
}

start_pre() {
	checkpath -d -m 0755 /run/nvidia-persistenced
}

stop_post() {
	rm -rf /run/nvidia-persistenced
}
