# etc/init.d/nvidia-powerd
#!/sbin/openrc-run

description="nvidia-powerd service"
command="/usr/bin/nvidia-powerd"
command_background="yes"
pidfile="/run/nvidia-powerd.pid"

depend() {
    need mountall
    after bootmisc
}
