# etc/init.d/nvidia-powerd
#!/sbin/openrc-run

description="Manages dynamic power states for NVIDIA GPUs"
command="/usr/bin/nvidia-powerd"
command_background="yes"
pidfile="/run/nvidia-powerd.pid"

depend() {
	# This daemon requires the NVIDIA kernel modules to be loaded.
	need modules
	# It should start after the basic system is up.
	after localmount
}
