# debian/postinst
#!/usr/bin/env bash

# -- Add service to runlevel.

rc-update add nvidia-powerd async

# -- Update services list.

rc-update -u
